# 01_jazzy-vscode-xfce-tigerVNC-novnc/supervisord.xfce.conf
[supervisord]
nodaemon=true
user=root   ; optional: silences the "running as root" warning

[program:vnc]
; Use gosu **-H** to set HOME to /home/${USER}, and call an absolute path.
command=/bin/bash -lc "exec gosu -H ${USER} /home/${USER}/.vnc/vnc_run.sh"
autostart=true
autorestart=true
stdout_logfile=/tmp/vnc.log
stderr_logfile=/tmp/vnc.err
startsecs=2

[program:novnc]
; If you use noVNC in this image, keep this; otherwise you can delete this block.
command=/bin/bash -lc "websockify --web=/usr/share/novnc ${NO_VNC_PORT:-6080} localhost:${VNC_PORT:-5901}"
autostart=true
autorestart=true
stdout_logfile=/tmp/novnc.log
stderr_logfile=/tmp/novnc.err
startsecs=2
